# profile.d/00-environment.sh
#!/usr/bin/env bash
################################
export    "SXHKD_SHELL=${SHELL}"
export    "TMPDIR=/tmp"
export    "EDITOR=nano"
################################
export    "MODE=1920x1080"
export    "SKEL=/etc/skel"
export    "SHARE=/var/http/shrd"
################################
