# profile.d/00-environment.sh
#!/usr/bin/env bash
##########################################
export    "SXHKD_SHELL=${SHELL}"
export    "TMPDIR=/tmp"
export    "EDITOR=nano"
##########################################
export    "NEED=nature,animals"
export    "MODE=1920x1080"
export    "SKEL=/etc/skel"
export    "IMG=/etc/skel/Picture/_.jpg"
export    "HTTP=/var/http"
export    "SHARE=${HTTP}/usrs/local/share"
##########################################
export    "THEMES=/usr/share/themes"
export    "ICONS=/usr/share/icons"
##########################################
