# profile.d/color.sh
#!/usr/bin/env bash
###############################
source    "colo.sh"
###############################
printf    "\e]P0${black}"
printf    "\e]P1${darkred}"
printf    "\e]P2${darkgreen}"
printf    "\e]P3${brown}"
printf    "\e]P4${darkblue}"
printf    "\e]P5${darkmagenta}"
printf    "\e]P6${darkcyan}"
printf    "\e]P7${lightgrey}"
printf    "\e]P8${darkgrey}"
printf    "\e]P9${red}"
printf    "\e]PA${green}"
printf    "\e]PB${yellow}"
printf    "\e]PC${blue}"
printf    "\e]PD${magenta}"
printf    "\e]PE${cyan}"
printf    "\e]PF${white}"
###############################
