# cron.daily/00-clean
#!/usr/bin/env bash
###################
source    "clea.sh"
###################
