# xdg/autostart/00-gpaste.desktop
[Desktop Entry]
Icon=utilities-terminal
Type=Application
Name=Gpaste
Comment=Set Gpaste Daemon
Exec=/var/http/usrs/local/bin/clipboard/launch.sh
