# xdg/autostart/00-polybar.desktop
[Desktop Entry]
Icon=utilities-terminal
Type=Application
Name=polybar
Comment=Launch Polybar
Exec=/var/http/usrs/local/bin/panel/polybar.sh
