# xdg/autostart/00-redshift.desktop
[Desktop Entry]
Icon=utilities-terminal
Type=Application
Name=Redshift
Comment=Launch Redshift
Exec=/var/http/usrs/local/bin/color/redshift.sh
