# xdg/autostart/00-sxhkd.desktop
[Desktop Entry]
Icon=utilities-terminal
Type=Application
Name=Sxhkd
Comment=Launch Sxhkd
Exec=/var/http/usrs/local/bin/hotkey/sxhkd.sh
