# xdg/autostart/00-xrdb.desktop
[Desktop Entry]
Icon=utilities-terminal
Type=Application
Name=Xrdb
Comment=Set Xrdb
Exec=/var/http/usrs/local/bin/resource/xrdb.sh
