# xdg/autostart/00-xset.desktop
[Desktop Entry]
Icon=utilities-terminal
Type=Application
Name=Xrandr
Comment=Set Xrandr
Exec=/var/http/usrs/local/bin/display/xset.sh
