# xdg/autostart/00-xrdb.desktop
[Desktop Entry]
Icon=utilities-terminal
Type=Application
Name=Xrdb
Comment=Set Xrdb
Exec=/var/http/srcs/bin/xrdb.sh
